# record-camera.service
[Unit]
Description=Record camera on startup

[Service]
ExecStart=/usr/bin/record-camera
After=multi-user.target
